# Server/pull-auto.service
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/libexec/pull-auto start
